# root/etc/s6-overlay/s6-rc.d/init-meloday-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

if [[ ! -f /config/config.yml ]]; then
    cp /config.yml /config/config.yml
fi

# permissions
lsiown -R abc:abc \
    /config